# templates/default/uwsgi.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "uWSGI"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

<% node.python.uwsgi.envs.each do |variable, value| %>
env <%= variable %>='<%= value %>'
<% end %>

exec <%= node.python.uwsgi.bin %> --ini /etc/uwsgi.ini

post-start script
  sleep 1 # giving the daemon a chance to create the socket file
  chmod -f 777 <%= node.python.uwsgi.socket %>
end

killed script
  rm -f <%= node.python.uwsgi.pidfile %>
  rm -f <%= node.python.uwsgi.socket %>
end script

post-stop script
  rm -f <%= node.python.uwsgi.pidfile %>
  rm -f <%= node.python.uwsgi.socket %>
end script
